# pkg/deb/bookworm/rtpengine-daemon.service
[Unit]
Description=NGCP RTP/media Proxy Daemon
After=network-online.target
After=remote-fs.target
Requires=network-online.target

[Service]
Type=notify
EnvironmentFile=/etc/default/rtpengine-daemon
PIDFile=/run/rtpengine-daemon.pid
ExecStartPre=/usr/sbin/rtpengine-iptables-setup start
ExecStart=/usr/sbin/rtpengine -f -E --no-log-timestamps --pidfile /run/rtpengine-daemon.pid --config-file /etc/rtpengine/rtpengine.conf
ExecStopPost=/usr/sbin/rtpengine-iptables-setup stop

[Install]
WantedBy=multi-user.target
