# el/rtpengine.service
[Unit]
Description=NGCP RTP/media Proxy Daemon
After=network-online.target
Requires=network-online.target

[Service]
Type=forking
User=ngcp-rtpengine
Group=daemon
Environment=CFGFILE=/etc/rtpengine/rtpengine-recording.conf
EnvironmentFile=/etc/sysconfig/rtpengine
Environment=CFGFILE=/etc/rtpengine/rtpengine.conf
PIDFile=/run/rtpengine.pid
ExecStartPre=/usr/sbin/ngcp-rtpengine-iptables-setup start
ExecStart=/usr/sbin/rtpengine --no-log-timestamps --pidfile=${PIDFile} --config-file=${CFGFILE} --table=${TABLE}
ExecStopPost=/usr/sbin/ngcp-rtpengine-iptables-setup stop

[Install]
WantedBy=multi-user.target
