# debian/ngcp-rtpengine-daemon.postrm
#!/bin/sh
# postrm script for ngcp-rtpengine-daemon

set -e

removal_wrapper() {
  # remove the init script only on ce systems, as the
  # the pro system handle it inside the monitoring/HA setup
  if ! [ -x "$(which ngcp-check_active 2>/dev/null)" ]; then
    update-rc.d ngcp-rtpengine-daemon remove >/dev/null
  fi
}

if [ "$1" = "purge" ] ; then
  removal_wrapper
fi

exit 0
