# debian/ngcp-rtpengine-kernel-dkms.prerm
#!/bin/sh

set -e

# make sure it's not running
if [ -x /etc/init.d/ngcp-rtpengine-daemon ] ; then
	invoke-rc.d ngcp-rtpengine-daemon stop || true
	rmmod "xt_RTPENGINE" 2>/dev/null || true
fi

#DEBHELPER#

exit 0
